# deploy/howick-agent-update.service
[Unit]
Description=howick-agent auto-update from GitHub Releases
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/howick-agent-update.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=howick-agent-update
