# deploy/howick-agent.service
[Unit]
Description=Howick Agent — USB Gadget Edge Agent (Pi Zero 2W)
Documentation=https://github.com/joeblew999/opcua-howick
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStartPre=/bin/sleep 5
ExecStart=/home/pi/howick-agent --config /home/pi/howick-agent.pi-zero.toml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=howick-agent

# Environment
Environment=RUST_LOG=howick_agent=info
Environment=RUST_BACKTRACE=1

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/jobs /mnt/usb_share

[Install]
WantedBy=multi-user.target
