# deploy/opcua-server.service
[Unit]
Description=OPC UA server for Howick FRAMA Machine (Pi 5)
Documentation=https://github.com/joeblew999/opcua-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/home/pi/opcua-server --config /home/pi/opcua-server.pi5.toml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opcua-server

# Environment
Environment=RUST_LOG=opcua_howick=info
Environment=RUST_BACKTRACE=1

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/jobs

[Install]
WantedBy=multi-user.target
